# recipes-core/inangoplug-rsc-init-arm-systemd/files/rsc-server.path
[Unit]
Description=Check flag file "VS Client enable"
After=puma_fs.service

[Path]
PathExists=/tmp/vsclient_enable

[Install]
WantedBy=multi-user.target
